# ops/systemd/lightld-signer.service
[Unit]
Description=Lightld local signer service
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/opt/lightld
ExecStart=/usr/bin/node --experimental-strip-types src/cli/run-local-signer-main.ts
Restart=always
RestartSec=5
User=lightld
EnvironmentFile=/etc/lightld-signer.env
ReadOnlyPaths=/opt/lightld/secrets

[Install]
WantedBy=multi-user.target
