# deploy/systemd/lightld-signer.service
[Unit]
Description=Lightld local signer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/lightld
Environment=LIGHTLD_ROOT=/opt/lightld
ExecStart=/usr/bin/bash /opt/lightld/start-signer.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
